# modules/development/files/systemd/system/docker.service.d/docker.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --ip-forward=true --storage-driver overlay
